# upstart/verbify-consumers-restart.conf
description "restart queue consumers"

task

start on verbify-restart or verbify-kill

script
    . /etc/default/verbify
    manage-consumers
end script
