# upstart/verbify-boot.conf
description "start verbify on boot"

task 
start on runlevel [2345]

exec initctl emit verbify-start