# upstart/verbify-consumers-start.conf
description "start up queue consumers"

task

start on runlevel [2345] or verbify-start

script
    . /etc/default/verbify
    manage-consumers
end script
